# server/scripts/atomic-server.service
[Unit]
After=network.target

[Service]
EnvironmentFile=/etc/environment
ExecStart=/usr/bin/atomic-server-fc.sh

[Install]
WantedBy=default.target